# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:redis]
command = redis-server
stdout_logfile=logs/redis_stdout.log
stderr_logfile=logs/redis_stderr.log

[program:mongodb]
command = mongod -f mongodb-dev.conf
stdout_logfile=logs/mongodb_stdout.log
stderr_logfile=logs/mongodb_stderr.log

[program:eventq]
command = paster serve --reload development.ini
stdout_logfile=logs/eventq_stdout.log
stderr_logfile=logs/eventq_stderr.log
